# systemd/labrad-node.service.d/labrad-node.conf
[Service]
Environment="LABRADHOST="
Environment="LABRADPASSWORD="
Environment="LABRADNODE="
Environment="LABRAD_TLS=off"

