# systemd/labrad-node.service.d/labrad-node.conf
[Service]
Environment="LABRADHOST=yesr12.colorado.edu"
Environment="LABRADPASSWORD=698gang"
Environment="LABRADNODE=yesr9"
Environment="LABRAD_TLS=off"

